# supervisord.conf
[supervisord]
nodaemon=true

[program:artisan]
command=/usr/local/bin/php /var/www/html/artisan serve --host=0.0.0.0 --port=8000
autostart=true
autorestart=true
stdout_logfile=/var/log/artisan.log
stderr_logfile=/var/log/artisan_err.log
user=laravel

[program:npm]
command=/usr/local/bin/npm run dev
autostart=true
autorestart=true
stdout_logfile=/var/log/npm.log
stderr_logfile=/var/log/npm_err.log
user=laravel
